# gather2-odk-importer/conf/supervisor.importer.conf
[program:gather2_importer]
command=uwsgi --ini /code/conf/uwsgi.ini
directory=/code
autostart=true
autorestart=true
stdout_logfile=/var/log/uwsgi.log
redirect_stderr=True
